# bash/crontab
# start
*/15 * * * * /sample/bash/test_cli.sh >> /var/log/cron.log 2>&1
# end
